# upload/docker/supervisor/asset_consumer_notify.ini
[program:asset_consumer_notify]
directory=/srv/app
command=/bin/sh -c '/srv/app/bin/console rabbitmq:consumer -vvv --messages=150 asset_consumer_notify || (sleep 10s && false)'
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=1000
