# infra/docker/php-fpm-base/worker/worker.ini
[program:messenger-consume-${WORKER_CHANNEL}]
directory=/srv/app
command=./bin/console messenger:consume ${WORKER_CHANNEL} -vvv --limit=${WORKER_LIMIT} --time-limit=${WORKER_TIME_LIMIT} --no-reset
user=1000
numprocs=1
startsecs=$SUPERVISOR_START_SECS
startretries=$SUPERVISOR_START_RETRIES
autorestart=true
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=$SUPERVISOR_STOP_WAIT_SECS
