# uploader/api/docker/supervisor/fast_events.ini
[program:fast_events]
directory=/srv/app
command=/bin/sh -c '/srv/app/bin/console rabbitmq:consumer -vvv --messages=150 --heartbeat-sender=pcntl fast_events || (sleep 10s && false)'
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=1000
